# config/openwrt-tvbox/etc/crontabs/root
// 52 02 * * * sleep 5 && touch /etc/banner && reboot   //每天2点52分路由器自动重启
3 3 12 12 * /usr/bin/nginx-util 'check_ssl'
0 5 * * * /etc/init.d/smartdns updatefiles
